# etc/profile.d/grep.sh
#!/usr/bin/env bash

export GREP_OPTIONS='--color=auto'
