# etc/profile.d/jena.sh
#!/usr/bin/env bash

export PATH=${PATH}:${JENA_HOME}/bin
