# etc/profile.d/latex.sh
#!/usr/bin/env bash

#alias tlmgr='/usr/share/texmf-dist/scripts/texlive/tlmgr.pl'
