# etc/profile.d/tmp.sh
#!/usr/bin/env bash

export TMPDIR=/var/tmp